# cli/innernet.target
[Unit]
Description=Target to allow restarting and stopping of all parts of innernet
